# youtube_audio_downloader/rootfs/etc/services.d/youtube-audio/run
#!/usr/bin/with-contenv bashio
exec uvicorn youtube_audio.main:create_app \
  --app-dir /app/backend \
  --factory \
  --host 0.0.0.0 \
  --port 8099 \
  --no-server-header
